# deploy/videomind.service
[Unit]
Description=VideoMind single-server web application
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__APP_USER__
Group=__APP_GROUP__
WorkingDirectory=__APP_DIR__/backend
EnvironmentFile=-__APP_DIR__/backend/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=__APP_DIR__/.venv/bin/python -m uvicorn main:app --host 127.0.0.1 --port 8000 --workers 1
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
KillMode=control-group
MemoryMax=3G
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
