# deploy/videomind-maintenance.service
[Unit]
Description=VideoMind backup and stale-file cleanup

[Service]
Type=oneshot
User=__APP_USER__
Group=__APP_GROUP__
WorkingDirectory=__APP_DIR__
EnvironmentFile=-__APP_DIR__/backend/.env
ExecStart=/usr/bin/bash __APP_DIR__/deploy/backup.sh
ExecStart=/usr/bin/bash __APP_DIR__/deploy/cleanup.sh
